# lib/ESPiLight/pilight/res/init/pilight.upstart
#!upstart

description "start pilight daemon"

start on started mountall and net-device-up
stop on runlevel [016]

expect fork

exec pilight-daemon